# test.timer
[Unit]
Description=Timer for test.service

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
Unit=test.service

[Install]
WantedBy=timers.target
